# install/13-install-sensors/lm_sensors.list
lm_sensors /etc/
lm_sensors /etc/conf.d/
lm_sensors /etc/conf.d/healthd
lm_sensors /etc/conf.d/sensord
lm_sensors /etc/rc.d/
lm_sensors /etc/rc.d/fancontrol
lm_sensors /etc/rc.d/healthd
lm_sensors /etc/rc.d/sensord
lm_sensors /etc/rc.d/sensors
lm_sensors /etc/sensors3.conf
lm_sensors /usr/
lm_sensors /usr/bin/
lm_sensors /usr/bin/sensors
lm_sensors /usr/bin/sensors-conf-convert
lm_sensors /usr/include/
lm_sensors /usr/include/sensors/
lm_sensors /usr/include/sensors/error.h
lm_sensors /usr/include/sensors/sensors.h
lm_sensors /usr/lib/
lm_sensors /usr/lib/libsensors.so
lm_sensors /usr/lib/libsensors.so.4
lm_sensors /usr/lib/libsensors.so.4.3.1
lm_sensors /usr/sbin/
lm_sensors /usr/sbin/fancontrol
lm_sensors /usr/sbin/healthd
lm_sensors /usr/sbin/isadump
lm_sensors /usr/sbin/isaset
lm_sensors /usr/sbin/pwmconfig
lm_sensors /usr/sbin/sensord
lm_sensors /usr/sbin/sensors-detect
lm_sensors /usr/share/
lm_sensors /usr/share/man/
lm_sensors /usr/share/man/man1/
lm_sensors /usr/share/man/man1/sensors.1.gz
lm_sensors /usr/share/man/man3/
lm_sensors /usr/share/man/man3/libsensors.3.gz
lm_sensors /usr/share/man/man5/
lm_sensors /usr/share/man/man5/sensors.conf.5.gz
lm_sensors /usr/share/man/man5/sensors3.conf.5.gz
lm_sensors /usr/share/man/man8/
lm_sensors /usr/share/man/man8/fancontrol.8.gz
lm_sensors /usr/share/man/man8/isadump.8.gz
lm_sensors /usr/share/man/man8/isaset.8.gz
lm_sensors /usr/share/man/man8/pwmconfig.8.gz
lm_sensors /usr/share/man/man8/sensord.8.gz
lm_sensors /usr/share/man/man8/sensors-detect.8.gz
